# .config/autostart/Udiskie.desktop
[Desktop Entry]
Type=Application
Name=Udiskie
Description=Automount helper
Exec=udiskie --tray
